# docker/supervisord.conf
[supervisord]
nodaemon=true
user=root
logfile=/dev/null
pidfile=/tmp/supervisord.pid

[program:python-model]
command=/usr/bin/python3 /app/app/Infrastructure/Prediction/Python/prediction_server.py
directory=/app
autostart=true
autorestart=true
startsecs=5
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
stderr_logfile=/dev/fd/2
stderr_logfile_maxbytes=0
environment=BRAIN_TUMOR_MODEL_PATH="/app/app/Infrastructure/Prediction/Python/best_model.pth",BRAIN_TUMOR_IMAGE_SIZE="224",BRAIN_TUMOR_CLASS_LABELS="glioma,meningioma,notumor,pituitary"

[program:laravel-web]
command=/usr/local/bin/php /app/artisan serve --host=0.0.0.0 --port=%(ENV_PORT)s
directory=/app
autostart=true
autorestart=true
startsecs=5
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
stderr_logfile=/dev/fd/2
stderr_logfile_maxbytes=0
